# crontab/ftp-puller-crontab
SHELL=/bin/bash
# Loading environment from predefined env file /root/.cronenv with format:
# varName1=value 
# varName2=value 
1 1 1 1 1 . $HOME/.cronenv; env | sort; /usr/local/bin/python /root/ftp-files-puller.py EXTRA_ARGS   >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job